# systemd/.config/systemd/user/wireplumber.service.d/override.conf
[Unit]
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStartPre=/usr/bin/sleep 2
